# supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=python app_gpt.py
autostart=true
autorestart=true
