# supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=python app_gpt.py
autostart=true
autorestart=true
stdout_logfile=/var/log/flask_stdout.log
stderr_logfile=/var/log/flask_stderr.log
